# com.greyhound.automation.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.greyhound.automation</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/orlandolee/greyhound_racing_collector/automation_scheduler.py</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/Users/orlandolee/greyhound_racing_collector</string>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <true/>
    
    <key>StandardOutPath</key>
    <string>/Users/orlandolee/greyhound_racing_collector/logs/automation/launchd_out.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/orlandolee/greyhound_racing_collector/logs/automation/launchd_err.log</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin</string>
        <key>PYTHONPATH</key>
        <string>/Users/orlandolee/greyhound_racing_collector</string>
    </dict>
    
    <key>ProcessType</key>
    <string>Background</string>
    
    <key>LowPriorityIO</key>
    <true/>
    
    <key>Nice</key>
    <integer>10</integer>
</dict>
</plist>
