# archive/com.greyhound.git-backup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.greyhound.git-backup</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/Users/orlandolee/greyhound_racing_collector/git_backup.sh</string>
        <string>scheduled</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/Users/orlandolee/greyhound_racing_collector</string>
    
    <key>StartCalendarInterval</key>
    <array>
        <!-- Daily backup at 2:00 AM -->
        <dict>
            <key>Hour</key>
            <integer>2</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <!-- Additional backup at 2:00 PM -->
        <dict>
            <key>Hour</key>
            <integer>14</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <!-- Evening backup at 10:00 PM -->
        <dict>
            <key>Hour</key>
            <integer>22</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
    </array>
    
    <key>StandardOutPath</key>
    <string>/Users/orlandolee/greyhound_racing_collector/logs/git_backup_stdout.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/orlandolee/greyhound_racing_collector/logs/git_backup_stderr.log</string>
    
    <key>RunAtLoad</key>
    <false/>
    
    <key>KeepAlive</key>
    <false/>
    
    <key>ProcessType</key>
    <string>Background</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>/Users/orlandolee</string>
    </dict>
</dict>
</plist>
